# .config/systemd/user/parcimonie.service
[Unit]
Description=GnuPG parcimonie key refresher

[Service]
ExecStart=/usr/bin/parcimonie.sh

[Install]
WantedBy=default.target
